# ftps/srcs/ftps.ini
[supervisord]
nodaemon=true

[program:vsftpd]
command=rc-service vsftpd start